# crontab
0 */2 * * * ~/Documents/Projects/bs-notify/bsnotify.py
